# perroGuardian.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/perroguardian.sh &

[Install]
WantedBy=multi-user.target
